# recorder/dist/recorder.service
[Unit]
Description=SwitchBot Recorder
Requires=redis-server.service

[Service]
Type=simple
WorkingDirectory=/opt/recorder
ExecStart=/opt/recorder/current/main
EnvironmentFile=/opt/recorder/environment
Restart=always

[Install]
WantedBy=multi-user.target
